# scripts/robotic-arm.service
[Unit]
Description=Robotic Arm Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/start-robotic-arm
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target